# etc/crontab
#minute	hour	mday	month	wday	who	command
*/1	*	*	*	*	root	/local/scripts/track
